# quadlets/open-webui.container
[Unit]
Description=Open WebUI
After=network-online.target ai-stack-pod.service ollama-rocm.service
Wants=network-online.target ai-stack-pod.service ollama-rocm.service

[Container]
ContainerName=open-webui
Image=ghcr.io/open-webui/open-webui:main
Pod=ai-stack.pod
Volume=%h/.local/share/open-webui:/app/backend/data:Z
Environment=OLLAMA_BASE_URL=http://ollama:11434

[Service]
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
